# crontab
*/2 * * * * docker run -v oli_data:/usr/local/argus/data -v oli_logs:/var/log/argus -e DISCOGS_TOKEN_KEY=discogs_token_oli -e TELEGRAM_CHAT_ID_KEY=telegram_chat_id_oli -e USER=oli -t argus:latest python main.py --task crawl_async
*/2 * * * * docker run -v pa_data:/usr/local/argus/data -v pa_logs:/var/log/argus -e DISCOGS_TOKEN_KEY=discogs_token_pa -e TELEGRAM_CHAT_ID_KEY=telegram_chat_id_pa -e USER=pa -t argus:latest python main.py --task crawl_async
34 */6 * * * docker run -v ash_data:/usr/local/argus/data -v ash_logs:/var/log/argus -e DISCOGS_TOKEN_KEY=discogs_token_ash -e TELEGRAM_CHAT_ID_KEY=telegram_chat_id_ash -e USER=ash -t argus:latest python main.py --task crawl_async
0 0 * * * docker system prune -f --filter "until=$((24))h"
